# deployment/systemd/airunner-headless.service
[Unit]
Description=AI Runner Headless Daemon
After=network.target

[Service]
Type=simple
User=__AIRUNNER_USER__
Group=__AIRUNNER_GROUP__
WorkingDirectory=__AIRUNNER_BUNDLE_ROOT__
Environment="AIRUNNER_HEADLESS=1"
Environment="AIRUNNER_BUNDLE_ROOT=__AIRUNNER_BUNDLE_ROOT__"
Environment="AIRUNNER_BASE_PATH=__AIRUNNER_DATA_DIR__"
Environment="AIRUNNER_DATA_DIR=__AIRUNNER_DATA_DIR__"
Environment="AIRUNNER_RUNTIME_ROOT=__AIRUNNER_RUNTIME_ROOT__"
Environment="AIRUNNER_RUNTIME_CONFIG_DIR=__AIRUNNER_RUNTIME_CONFIG_DIR__"
Environment="AIRUNNER_RUNTIME_LOG_DIR=__AIRUNNER_RUNTIME_LOG_DIR__"
Environment="AIRUNNER_RUNTIME_SOCKET_DIR=__AIRUNNER_RUNTIME_SOCKET_DIR__"
Environment="AIRUNNER_CACHE_DIR=__AIRUNNER_CACHE_DIR__"
Environment="AIRUNNER_MODEL_DIR=__AIRUNNER_MODEL_DIR__"
Environment="AIRUNNER_PYTHON=__AIRUNNER_PYTHON__"
Environment="AIRUNNER_DAEMON_CONFIG=__AIRUNNER_DAEMON_CONFIG__"
Environment="AIRUNNER_RUNTIME_BIND_HOST=127.0.0.1"
Environment="AIRUNNER_HTTP_HOST=127.0.0.1"
Environment="AIRUNNER_HTTP_PORT=8080"
Environment="AIRUNNER_LLM_ON=1"
Environment="QT_QPA_PLATFORM=offscreen"
Environment="QT_LOGGING_RULES=*.debug=false;qt.qpa.*=false"
Environment="AIRUNNER_LOG_LEVEL=INFO"
Environment="DEV_ENV=0"
Environment="PATH=__AIRUNNER_BIN_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

ExecStartPre=/bin/mkdir -p __AIRUNNER_RUNTIME_ROOT__ __AIRUNNER_RUNTIME_CONFIG_DIR__ __AIRUNNER_RUNTIME_LOG_DIR__ __AIRUNNER_RUNTIME_SOCKET_DIR__ __AIRUNNER_CACHE_DIR__ __AIRUNNER_MODEL_DIR__
ExecStartPre=/bin/chmod 700 __AIRUNNER_RUNTIME_ROOT__ __AIRUNNER_RUNTIME_CONFIG_DIR__ __AIRUNNER_RUNTIME_LOG_DIR__ __AIRUNNER_RUNTIME_SOCKET_DIR__ __AIRUNNER_CACHE_DIR__ __AIRUNNER_MODEL_DIR__

ExecStart=__AIRUNNER_PYTHON__ -m airunner.services.daemon --config __AIRUNNER_DAEMON_CONFIG__
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM
# Python logging handles file output, systemd just captures for journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=airunner-headless

# Performance and security settings
LimitNOFILE=65536
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=__AIRUNNER_DATA_DIR__
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
